# mail/rootfs/etc/services.d/php-fpm/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Mailserver
# Runs PHP-FPM
# ==============================================================================
bashio::log.info "Starting PHP-FPM..."

export host
export password
export port
export username

host=$(bashio::services "mysql" "host")
password=$(bashio::services "mysql" "password")
port=$(bashio::services "mysql" "port")
username=$(bashio::services "mysql" "username")

exec php-fpm84 --nodaemonize
