# mail/rootfs/etc/services.d/dovecot/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Mailserver
# Runs the Dovecot server
# ==============================================================================
bashio::log.info "Starting Dovecot...."

exec /usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf
