# mail/rootfs/etc/services.d/postfix/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant App: Mailserver
# Runs the Nginx Postfix server
# ==============================================================================

bashio::log.info "Starting Postfix...."
exec /usr/libexec/postfix/master -d
