# mail/rootfs/etc/s6-overlay/s6-rc.d/dovecot/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant App: Mailserver
# Runs the Dovecot server
# ==============================================================================
bashio::log.info "Starting Dovecot...."

exec /usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf
